# snoqtt-alpine-pcap/require/libs/super_snort.conf
[inet_http_server]
port = 127.0.0.1:9001

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = http://127.0.0.1:9001

[program:snort_parser]
command=python3 /root/libs/snort_parser.py

[program:snort_api]
command=python3 /root/libs/snort_api.py

[program:snort]
command=/usr/local/bin/snort -i %(ENV_NETINT)s -A unsock -N --daq pcap --daq-mode read-file -r /root/pcap/wednesday.pcap -c /etc/snort/snort.conf
#command=/usr/local/bin/snort -i %(ENV_NETINT)s -A unsock -c /etc/snort/snort.conf
#command=/usr/bin/snort -r /snort_cap/remotecapture.fifo -A unsock -c /etc/snort/snort.conf
autorestart=true
